# ci/ansible/roles/deploy/templates/tlm.socket
[Socket]
ListenStream=/run/tlm.sock
SocketUser=www-data
SocketMode=600

[Install]
WantedBy=sockets.target